# deploy/systemd/miniappupgr-deploy.service
[Unit]
Description=Deploy miniappupgr from GitHub
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
Environment=REPO_DIR=/opt/miniappupgr
Environment=DEPLOY_BRANCH=main
ExecStart=/bin/bash /opt/miniappupgr/deploy/deploy.sh
